# config/open-vm-tools_2/open-vm-tools.inc
<?php

function open_vm_tools_deinstall() {
	exec("rm /usr/local/etc/rc.d/vmware*");
	unlink_if_exists("/boot/kernel/vmblock.ko");
	unlink_if_exists("/boot/kernel/vmhgfs.ko");
	unlink_if_exists("/boot/kernel/vmmemctl.ko");
	unlink_if_exists("/boot/kernel/vmxnet.ko");
}

function open_vm_tools_install() {

	// clean up old .ko files if they exist
	unlink_if_exists("/boot/kernel/vmblock.ko");
	unlink_if_exists("/boot/kernel/vmhgfs.ko");
	unlink_if_exists("/boot/kernel/vmmemctl.ko");
	unlink_if_exists("/boot/kernel/vmxnet.ko");

	$openvmtools_path = "/usr/pbi/open-vm-tools-" . php_uname("m");
	
	exec("cp $openvmtools_path/local/lib/vmware-tools/modules/drivers/vmmemctl /boot/kernel/");
	
	$vmware_tools_conf = <<<EOF
[powerops]
poweron-script=
resume-script=
suspend-script=
poweroff-script=

EOF;

	$vmware_guestd = <<<EOF
#!/bin/sh
#
# This file was automatically generated 
# by the pfSense package manager.
# 
# Do not edit this file.  Edit 
# /usr/local/pkg/open-vm-tools.inc instead.
#
# PROVIDE: vmware-guestd
# REQUIRE: DAEMON
# BEFORE: LOGIN

. /etc/rc.subr

# Global
checkvm_cmd="/usr/local/bin/vmware-checkvm > /dev/null"

# VMware guest daemon
name="vmware_guestd"
rcvar="\${name}_enable"
start_precmd="\${checkvm_cmd}"
unset start_cmd
stop_precmd="\${checkvm_cmd}"
unset stop_cmd
command="/usr/local/bin/vmtoolsd"
command_args="-c {$openvmtools_path}/share/vmware-tools/tools.conf -p {$openvmtools_path}/local/lib/open-vm-tools/plugins/vmsvc"
pidfile="/var/run/\${name}.pid"

load_rc_config \$name
vmware_guestd_enable="YES"
vmware_guestd_flags="--background \${pidfile}"
run_rc_command "\$1"

EOF;

	$vmware_kmod = <<<EOF
#!/bin/sh
#
# This file was automatically generated 
# by the pfSense package manager.
# 
# Do not edit this file.  Edit 
# /usr/local/pkg/open-vm-tools.inc instead.
#
# PROVIDE: vmware-kmod
# REQUIRE: FILESYSTEMS
# BEFORE: netif

. /etc/rc.subr

# Global
checkvm_cmd="/usr/local/bin/vmware-checkvm > /dev/null"

# Functions
vmware_guest_vmmemctl_start()
{
	echo 'Loading vmmemctl kernel module.'
	kldload vmmemctl.ko >/dev/null 2>&1
}

# VMware kernel module: vmmemctl
name="vmware_guest_vmmemctl"
rcvar="\${name}_enable"
start_precmd="\${checkvm_cmd}"
start_cmd="vmware_guest_vmmemctl_start"
stop_precmd="\${checkvm_cmd}"
stop_cmd=":"

load_rc_config \$name
vmware_guest_vmmemctl_enable="YES"
run_rc_command "\$1"

EOF;

	// Write out conf files.
	$fd = fopen("/usr/local/etc/rc.d/vmware-guestd.sh", "w");
	if (!$fd) {
		log_error("Could not open /usr/local/etc/rc.d/vmware-guestd.sh for writing");
		die("Could not open /usr/local/etc/rc.d/vmware-guestd.sh for writing");
	}
	fwrite($fd, $vmware_guestd);
	fclose($fd);
	$fd = fopen("/usr/local/etc/rc.d/vmware-kmod.sh", "w");
	if (!$fd) {
		log_error("Could not open /usr/local/etc/rc.d/vmware-kmod.sh for writing");
		die("Could not open /usr/local/etc/rc.d/vmware-kmod.sh for writing");
	}
	fwrite($fd, $vmware_kmod);
	fclose($fd);
	$fd = fopen("$openvmtools_path/local/share/vmware-tools/tools.conf", "w");
	if (!$fd) {
		log_error("Could not open $openvmtools_path/local/share/vmware-tools/tools.conf for writing");
		die("Could not open $openvmtools_path/local/share/vmware-tools/tools.conf for writing");
	}
	fwrite($fd, $vmware_tools_conf);
	fclose($fd);

	// Remove non used files
	unlink_if_exists("/usr/local/etc/rc.d/vmware-kmod");
	unlink_if_exists("/usr/local/etc/rc.d/vmware-guestd");

	// Make sure files are executable.
	exec("chmod a+rx /usr/local/etc/rc.d/*.sh");

}

?>